# usr/lib/systemd/system/android-mount.service.d/40-slot-b.conf
[Service]
ExecStartPre=-/usr/libexec/krypton-bootslot
