# usr/lib/systemd/system/ofono.service.d/30-mtk.conf
[Service]
Environment=OFONO_BINDER_PLUGIN_MTK_RSSI=true
