# usr/lib/systemd/system/phosh.service.d/45-brightness.conf
[Service]
ExecStartPre=-/usr/libexec/krypton-brightness
